# docs/examples/ganeti.cron
# restart failed instances
*/5 * * * * root /usr/local/sbin/ganeti-watcher
