# doc/examples/systemd/ganeti-rapi.service.in
[Unit]
Description = Ganeti RAPI daemon (rapi)
Documentation = man:ganeti-rapi(8)
Requires = ganeti-common.service
Requisite = ganeti-luxid.service
After = ganeti-common.service
PartOf = ganeti-master.target
ConditionPathExists = @LOCALSTATEDIR@/lib/ganeti/rapi.pem

[Service]
Type = simple
User = @GNTRAPIUSER@
Group = @GNTRAPIGROUP@
ExecStart = @SBINDIR@/ganeti-rapi -f
SuccessExitStatus = 0 11
Restart = on-failure

[Install]
WantedBy = ganeti-master.target ganeti.target
