# doc/examples/systemd/ganeti-common.service.in
[Unit]
Description = Ganeti one-off setup

[Service]
Type = oneshot
ExecStart = @PKGLIBDIR@/ensure-dirs
