# doc/examples/systemd/ganeti-kvmd.service.in
[Unit]
Description = Ganeti KVM daemon (kvmd)
Documentation = man:ganeti-kvmd(8)
Requires = ganeti-common.service
After = ganeti-common.service
PartOf = ganeti-noded.target

[Service]
Type = simple
Group = @GNTDAEMONSGROUP@
ExecStart = @SBINDIR@/ganeti-kvmd -f
Restart = on-failure

[Install]
WantedBy = ganeti-node.target ganeti.target
